# roles/configure/templates/ansible-pull@.service.j2
[Unit]
Description=Run ansible-pull %i
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ansible-pull-run %i
