# roles/configure/templates/warp-check.service.j2
[Unit]
Description=Check WARP CLI status and restart if not connected
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-warp-status.sh